# server/deploy/quiz-consumer.service
[Unit]
Description=Quiz Answers Consumer
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/quiz/server
Environment=NODE_ENV=production
EnvironmentFile=/opt/quiz/server/.env
ExecStart=/usr/bin/node /opt/quiz/server/dist/consumer.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
